# scripts/bingwallpaper.timer
[Unit]
Description=Bing Wallpaper Update

[Timer]
OnBootSec=5m
OnCalendar=*-*-* 05:00:00
Persistent=true

[Install]
WantedBy=timers.target
