# scripts/bingwallpaper.service
[Unit]
Description=Update Bing Wallpaper

[Service]
Type=oneshot
ExecStart=/home/kirkham/dotfiles/scripts/get_bing_image.sh
